# home/dot_config/autostart/fcitx5.desktop
# Managed by chezmoi. Do not edit directly. Update via chezmoi, using OS/host-specific partitioning when needed.
[Desktop Entry]
Type=Application
Exec=fcitx5
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[ja_JP]=fcitx5
Name=fcitx5
Comment[ja_JP]=fcitx5
Comment=fcitx5
