# home/dot_config/autostart/flatpak-1.desktop
# Managed by chezmoi. Do not edit directly. Update via chezmoi, using OS/host-specific partitioning when needed.
[Desktop Entry]
Type=Application
Exec=flatpak run io.github.jeffshee.Hidamari
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[ja_JP]=Hidamari
Name=Hidamari
Comment[ja_JP]=Live Desktop
Comment=Live Desktop
