# home/dot_config/autostart/flatpak.desktop
# Managed by chezmoi. Do not edit directly. Update via chezmoi, using OS/host-specific partitioning when needed.
[Desktop Entry]
Type=Application
Exec=flatpak run com.rustdesk.RustDesk
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[ja_JP]=RustDesk
Name=RustDesk
Comment[ja_JP]=
Comment=
