# home/dot_config/autostart/steam.desktop
[Desktop Entry]
Name=Steam
Comment=steam
Exec=steam
Type=Application
Terminal=false
Hidden=false

