# home/dot_config/autostart/steam.desktop
# Managed by chezmoi. Do not edit directly. Update via chezmoi, using OS/host-specific partitioning when needed.
[Desktop Entry]
Name=Steam
Comment=steam
Exec=steam
Type=Application
Terminal=false
Hidden=false

