# home/dot_config/systemd/user/rgb-load.service
# Managed by chezmoi. Do not edit directly. Update via chezmoi, using OS/host-specific partitioning when needed.
[Unit]
Description=Set RGB based on CPU load (one-shot)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rgb_load_once.sh
