# home/dot_config/systemd/user/rgb-load.timer
# Managed by chezmoi. Do not edit directly. Update via chezmoi, using OS/host-specific partitioning when needed.
[Unit]
Description=Update RGB periodically

[Timer]
OnBootSec=10
OnUnitActiveSec=1
AccuracySec=1s

[Install]
WantedBy=timers.target
